# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/backend/laravel/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker.log
stopasgroup=true

[program:laravel-scheduler]
process_name=%(program_name)s
command=php /app/backend/laravel/artisan schedule:work
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stopasgroup=true

[program:nginx]
process_name=%(program_name)s
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stopasgroup=true

[program:php-fpm]
process_name=%(program_name)s
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/php-fpm.log
stopasgroup=true

[program:fastapi]
process_name=%(program_name)s
command=uvicorn main:app --host 0.0.0.0 --port 8001 --workers 4
directory=/app/backend/fastapi
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/fastapi.log
environment=PYTHONUNBUFFERED=1,DATABASE_URL=%(ENV_DATABASE_URL)s
stopasgroup=true

[group:whale]
programs=laravel-worker,laravel-scheduler,nginx,php-fpm,fastapi

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface